# standard-armv7l/etc/profile.d/buildsystem.sh
export SUSE_IGNORED_RPATHS=/etc/ld.so.conf
